# socket-proxy/rootfs/etc/services.d/socket-proxy/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting Docker Socket Proxy..."

# Convert bashio bool (true/false) to 0/1 for haproxy env() checks
bool2int() {
    if bashio::config.true "$1"; then echo 1; else echo 0; fi
}

# Read options and export as environment variables
ALLOW_START=$(bool2int 'ALLOW_START')
ALLOW_STOP=$(bool2int 'ALLOW_STOP')
ALLOW_RESTARTS=$(bool2int 'ALLOW_RESTARTS')
AUTH=$(bool2int 'AUTH')
BUILD=$(bool2int 'BUILD')
COMMIT=$(bool2int 'COMMIT')
CONFIGS=$(bool2int 'CONFIGS')
CONTAINERS=$(bool2int 'CONTAINERS')
DISABLE_IPV6=$(bool2int 'DISABLE_IPV6')
DISTRIBUTION=$(bool2int 'DISTRIBUTION')
EVENTS=$(bool2int 'EVENTS')
EXEC=$(bool2int 'EXEC')
GRPC=$(bool2int 'GRPC')
IMAGES=$(bool2int 'IMAGES')
INFO=$(bool2int 'INFO')
LOG_LEVEL=$(bashio::config 'LOG_LEVEL')
NETWORKS=$(bool2int 'NETWORKS')
NODES=$(bool2int 'NODES')
PING=$(bool2int 'PING')
PLUGINS=$(bool2int 'PLUGINS')
POST=$(bool2int 'POST')
SECRETS=$(bool2int 'SECRETS')
SERVICES=$(bool2int 'SERVICES')
SESSION=$(bool2int 'SESSION')
SWARM=$(bool2int 'SWARM')
SYSTEM=$(bool2int 'SYSTEM')
TASKS=$(bool2int 'TASKS')
VERSION=$(bool2int 'VERSION')
VOLUMES=$(bool2int 'VOLUMES')

export ALLOW_START ALLOW_STOP ALLOW_RESTARTS AUTH BUILD COMMIT CONFIGS
export CONTAINERS DISABLE_IPV6 DISTRIBUTION EVENTS EXEC GRPC IMAGES INFO
export LOG_LEVEL NETWORKS NODES PING PLUGINS POST SECRETS SERVICES
export SESSION SWARM SYSTEM TASKS VERSION VOLUMES

PORT=$(bashio::config 'port')
SOCKET_PATH=/var/run/docker.sock
export SOCKET_PATH

# Determine bind protocol
if [ "${DISABLE_IPV6}" = "1" ]; then
    BIND_PROTO=":${PORT}"
else
    BIND_PROTO="ipv6@:${PORT}"
fi

# Prepare haproxy config
mkdir -p /run/haproxy
sed "s|@@BIND_PROTO@@|${BIND_PROTO}|g" /templates/haproxy.cfg > /run/haproxy/haproxy.cfg

bashio::log.info "Proxy listening on port ${PORT}"
bashio::log.info "IPv6: $([ "${DISABLE_IPV6}" = "1" ] && echo 'disabled' || echo 'enabled')"

exec haproxy -f /run/haproxy/haproxy.cfg -W -db
