# socket-proxy/rootfs/etc/services.d/socket-proxy/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

# Startup banner
bashio::log.info ""
bashio::log.info "  +------------------------------------+"
bashio::log.info "  |  ____             _                |"
bashio::log.info "  | |  _ \  ___   ___| | _____ _ __    |"
bashio::log.info "  | | | | |/ _ \ / __| |/ / _ \ '__|   |"
bashio::log.info "  | | |_| | (_) | (__|   <  __/ |      |"
bashio::log.info "  | |____/ \___/ \___|_|\_\___|_|      |"
bashio::log.info "  |                                    |"
bashio::log.info "  |  ____             _        _       |"
bashio::log.info "  | / ___|  ___   ___| | _____| |_     |"
bashio::log.info "  | \___ \ / _ \ / __| |/ / _ \ __|    |"
bashio::log.info "  |  ___) | (_) | (__|   <  __/ |_     |"
bashio::log.info "  | |____/ \___/ \___|_|\_\___|\__|    |"
bashio::log.info "  |                                    |"
bashio::log.info "  |  ____                              |"
bashio::log.info "  | |  _ \ _ __ _____  ___   _         |"
bashio::log.info "  | | |_) | '__/ _ \ \/ / | | |        |"
bashio::log.info "  | |  __/| | | (_) >  <| |_| |        |"
bashio::log.info "  | |_|   |_|  \___/_/\_\\\__,|        |"
bashio::log.info "  |                      |___/         |"
bashio::log.info "  +------------------------------------+"
bashio::log.info ""

# Version and config info
ADDON_VERSION="1.0.1"
UPSTREAM_PIN="linuxserver/docker-socket-proxy:3.2.13-r0-ls70"
bashio::log.info "Version:   ${ADDON_VERSION}"
bashio::log.info "Upstream:  ${UPSTREAM_PIN}"

# Convert bashio bool (true/false) to 0/1 for haproxy env() checks
bool2int() {
    if bashio::config.true "$1"; then echo 1; else echo 0; fi
}

# Read options and export as environment variables
ALLOW_START=$(bool2int 'ALLOW_START')
ALLOW_STOP=$(bool2int 'ALLOW_STOP')
ALLOW_RESTARTS=$(bool2int 'ALLOW_RESTARTS')
AUTH=$(bool2int 'AUTH')
BUILD=$(bool2int 'BUILD')
COMMIT=$(bool2int 'COMMIT')
CONFIGS=$(bool2int 'CONFIGS')
CONTAINERS=$(bool2int 'CONTAINERS')
DISABLE_IPV6=$(bool2int 'DISABLE_IPV6')
DISTRIBUTION=$(bool2int 'DISTRIBUTION')
EVENTS=$(bool2int 'EVENTS')
EXEC=$(bool2int 'EXEC')
GRPC=$(bool2int 'GRPC')
IMAGES=$(bool2int 'IMAGES')
INFO=$(bool2int 'INFO')
LOG_LEVEL=$(bashio::config 'LOG_LEVEL')
NETWORKS=$(bool2int 'NETWORKS')
NODES=$(bool2int 'NODES')
PING=$(bool2int 'PING')
PLUGINS=$(bool2int 'PLUGINS')
POST=$(bool2int 'POST')
SECRETS=$(bool2int 'SECRETS')
SERVICES=$(bool2int 'SERVICES')
SESSION=$(bool2int 'SESSION')
SWARM=$(bool2int 'SWARM')
SYSTEM=$(bool2int 'SYSTEM')
TASKS=$(bool2int 'TASKS')
VERSION=$(bool2int 'VERSION')
VOLUMES=$(bool2int 'VOLUMES')

export ALLOW_START ALLOW_STOP ALLOW_RESTARTS AUTH BUILD COMMIT CONFIGS
export CONTAINERS DISABLE_IPV6 DISTRIBUTION EVENTS EXEC GRPC IMAGES INFO
export LOG_LEVEL NETWORKS NODES PING PLUGINS POST SECRETS SERVICES
export SESSION SWARM SYSTEM TASKS VERSION VOLUMES

PORT=$(bashio::config 'port')

# Detect Docker socket location
if [ -S /var/run/docker.sock ]; then
    SOCKET_PATH=/var/run/docker.sock
elif [ -S /run/docker.sock ]; then
    SOCKET_PATH=/run/docker.sock
else
    bashio::log.fatal "Docker socket not found! Ensure 'docker_api: true' is set in config.yaml."
    exit 1
fi
export SOCKET_PATH

# Determine bind protocol
if [ "${DISABLE_IPV6}" = "1" ]; then
    BIND_PROTO=":${PORT}"
else
    BIND_PROTO="ipv6@:${PORT}"
fi

# Prepare haproxy config
mkdir -p /run/haproxy
sed "s|@@BIND_PROTO@@|${BIND_PROTO}|g" /templates/haproxy.cfg > /run/haproxy/haproxy.cfg

# Build list of enabled endpoints
ENABLED=()
for OPT in AUTH BUILD COMMIT CONFIGS CONTAINERS DISTRIBUTION EVENTS EXEC \
           GRPC IMAGES INFO NETWORKS NODES PING PLUGINS SECRETS SERVICES \
           SESSION SWARM SYSTEM TASKS VERSION VOLUMES; do
    if [ "${!OPT}" = "1" ]; then
        ENABLED+=("${OPT}")
    fi
done

# Build list of enabled write operations
WRITE_OPS=()
if [ "${POST}" = "1" ]; then WRITE_OPS+=("POST"); fi
if [ "${ALLOW_START}" = "1" ]; then WRITE_OPS+=("ALLOW_START"); fi
if [ "${ALLOW_STOP}" = "1" ]; then WRITE_OPS+=("ALLOW_STOP"); fi
if [ "${ALLOW_RESTARTS}" = "1" ]; then WRITE_OPS+=("ALLOW_RESTARTS"); fi

bashio::log.info "Port:      ${PORT}"
bashio::log.info "IPv6:      $([ "${DISABLE_IPV6}" = "1" ] && echo 'disabled' || echo 'enabled')"
bashio::log.info "Log level: ${LOG_LEVEL}"
bashio::log.info "Enabled:   ${ENABLED[*]:-none}"
if [ ${#WRITE_OPS[@]} -gt 0 ]; then
    bashio::log.warning "Write ops: ${WRITE_OPS[*]}"
fi
bashio::log.info ""

exec haproxy -f /run/haproxy/haproxy.cfg -W -db
